# docker/supervisor.conf
[program:vite]
command=npm run dev
directory=/var/www/html
autostart=true
autorestart=true
user=sail
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[program:reverb]
# Reverb 已被移除，此配置禁用該進程
autostart=false
autorestart=false
